# crontab
crontab -e

@reboot cd /home/pi/rpi-vocal-fx-unit && ./launch
